# docker/workspace/crontab/laravel
* * * * * laravel /usr/bin/php ${APP_CODE_PATH}/artisan schedule:run >> /dev/null 2>&1
